# package/obsbugzilla-sink.timer
[Timer]
OnCalendar=*:2,12,22,32,42,52
RandomizedDelaySec=15
Persistent=true

[Install]
WantedBy=timers.target
